# app_core/deploy/gunicorn.service
# deploy/gunicorn.service
# Minimal example systemd unit for running Gunicorn.
#
# Operators MUST replace placeholders below (paths, user, environment) to suit their environment.
#
# - ExecStart uses a conservative virtualenv path placeholder: /opt/venv/bin/gunicorn
# - Gunicorn config file path placeholder: /srv/smartlink/gunicorn.conf.py
#
# Example Environment lines are included for DATABASE_URL and WORKER_COUNT.
#
# Notes:
#  - TimeoutStopSec is tuned to allow graceful workers to exit within graceful_timeout configured in gunicorn.conf.py.
#  - Restart=on-failure is conservative; change to 'always' if you prefer supervisor-style restart on any exit.
#
[Unit]
Description=Gunicorn instance for Smartlink (example)
After=network.target

[Service]
Type=simple
# Recommended: run as a dedicated, unprivileged user
# User=www-data
# Group=www-data

# Working directory for your application (replace with real path)
# e.g. /srv/smartlink
WorkingDirectory=/srv/smartlink

# Path to the gunicorn binary from your virtualenv - replace as appropriate
ExecStart=/opt/venv/bin/gunicorn wsgi.prod_wsgi:app --config /srv/smartlink/gunicorn.conf.py

# Example environment variables (operators should set these in an environment file or here)
# Environment="DATABASE_URL=postgresql://user:pass@db.example:5432/smartlink"
# Environment="WORKER_COUNT=5"
# Optionally use EnvironmentFile=/etc/default/gunicorn

# Restart policy - conservative
Restart=on-failure
RestartSec=5

# Allow workers time to shutdown gracefully. Matches or slightly exceeds Gunicorn graceful_timeout.
TimeoutStopSec=45

# File descriptor limit and process accounting (optional)
# LimitNOFILE=65536

# Ensure journald captures stdout/stderr for container-friendly logging if ExecStart uses '-' logs
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
--- END FILE ---